# desktop-env/.config/systemd/user/autotz.service
[Unit]
Description=Auto-set timezone from geoclue location

[Service]
Type=oneshot
ExecStart=%h/.local/bin/autotz
