# desktop-env/.config/systemd/user/battery-alert.service
[Unit]
Description=Battery level alert and auto-suspend

[Service]
Type=oneshot
ExecStart=%h/.config/battery-alert/battery-alert.sh
