# desktop-env/.config/systemd/user/micro-exercise.timer
[Unit]
Description=Periodic micro exercise reminder (every 15 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
