# desktop-env/.config/systemd/user/micro-exercise.service
[Unit]
Description=Micro exercise reminder

[Service]
Type=oneshot
ExecStart=%h/.config/micro-exercise/micro-exercise.sh
