# system/etc/udev/rules.d/99-lid-suspend.rules
ACTION=="change", SUBSYSTEM=="button", KERNEL=="LID", ATTR{state}=="closed", RUN+="/usr/bin/systemctl --no-block start worth-suspend.service"
