# system/etc/systemd/system/worth-suspend.service
[Unit]
Description=Worth — lock then suspend

[Service]
Type=oneshot
ExecStart=/usr/local/bin/worth-suspend
TimeoutSec=30
